# tomenet-compiler.service
# container-tomenet.service
# autogenerated by Podman 2.1.1
# Mon Dec  7 09:55:06 CET 2020

[Unit]
Description=Podman container-tomenet.service
Documentation=man:podman-generate-systemd(1)
Wants=network.target
After=network-online.target

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=on-failure
ExecStart=/home/tokariew/.local/bin/make-tomenet-server
Type=oneshot

[Install]
WantedBy=multi-user.target default.target
